# devsend.service
[Unit]
Description=DevSend Email Automation Service
After=network.target

[Service]
Type=simple
User=tom
Group=prod
WorkingDirectory=/home/devsend
Environment="PATH=/home/devsend/venv/bin"
Environment="PYTHONPATH=/home/devsend"
ExecStart=/home/devsend/venv/bin/python3 -m devsend.main
Restart=always
RestartSec=10

# Logging (uses systemd journal - view with: journalctl -u devsend -f)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
