# devsend.service
[Unit]
Description=DevSend Email Automation Service
After=network.target

[Service]
Type=simple
User=tom
Group=prod
WorkingDirectory=/home/devsend
Environment="PATH=/home/devsend/venv/bin"
ExecStart=/home/devsend/venv/bin/uvicorn devsend.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/home/devsend/logs/devsend.log
StandardError=append:/home/devsend/logs/devsend-error.log

[Install]
WantedBy=multi-user.target
